# heartbleed/conf/supervisor/conf.d/bootstrap.conf
[program:bootstrap]
stdout_syslog=false
stderr_syslog=false
autorestart=false
startsecs=0
command=/etc/supervisor/scripts/bootstrap.sh
