# heartbleed/conf/supervisor/conf.d/nginx.conf
[program:nginx]
stdout_syslog=true
stderr_syslog=true
autostart=false
command=/usr/sbin/nginx
